# upstart/pgblackboard.conf
start on runlevel [2345]
stop on runlevel [!2345]
setuid pgblackboard
setgid pgblackboard
respawn
exec /usr/bin/env python3 -m pgblackboard
