# upstart/pgblackboard.conf
start on runlevel [2345]
stop on runlevel [!2345]
setuid pgblackboard
setgid pgblackboard
exec /usr/bin/env python3 -m pgblackboard.server --conf /etc/pgblackboard.conf
